# mycelium-managed/relay/deploy/control-plane.service
[Unit]
Description=Mycelium managed control-plane (provision + frps auth plugin)
After=network-online.target
Wants=network-online.target

[Service]
User=mycelium
Group=mycelium
WorkingDirectory=/opt/mycelium-managed
EnvironmentFile=/etc/mycelium/control-plane.env
ExecStart=/usr/bin/node src/server.js
Restart=on-failure
RestartSec=2
# registry.db lives in /var/lib/mycelium (created + owned by systemd).
StateDirectory=mycelium
ReadWritePaths=/var/lib/mycelium
# Hardening — the control-plane holds the DNS token + the handle registry.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_INET AF_INET6
RestrictNamespaces=true
LockPersonality=true

[Install]
WantedBy=multi-user.target
